# systemd/claude-mailbox.service
[Unit]
Description=Claude Mailbox Telegram Bot
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/claude-mailbox
ExecStart=/usr/bin/python3 /home/ubuntu/claude-mailbox/telegram_bot.py
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
